# .config/systemd/user/hapi-frp.target
[Unit]
Description=HAPI Hub with FRP Client
Requires=hapi-hub.service frpc-hapi.service
After=hapi-hub.service frpc-hapi.service
