# .config/systemd/user/waypaper-random.service
[Unit]
Description=Random Wallpaper via Waypaper

[Service]
ExecStart=%h/.config/waypaper/scripts/random-wallpaper-with-overview-sync.sh
