# .config/systemd/user/waypaper-random.service
[Unit]
Description=Random Wallpaper via Waypaper

[Service]
ExecStart=/usr/bin/waypaper --random
