# .config/systemd/user/waypaper-random.timer
[Unit]
Description=Timer to change wallpaper every 5 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=1800

[Install]
WantedBy=default.target
